# .config/systemd/user/checkupdates.service
[Unit]
Description=Checks for pacman updates.

[Service]
Type=oneshot
ExecStart=/home/$USER/.config/sh/checkupdates.sh